# ddns.service
[Unit]
Description=Cloudflare DDNS Service

[Service]
ExecStart=/usr/bin/node /home/serverkb/dynamicDNS/index.js
